# crontab
25	*	*	*	*	root	/usr/local/bin/tor-exit-nodes.sh
